# ansible/roles/supervisord/templates0/etc/supervisor/conf.d/template.conf
[program:{{ app_name }}]
command=/app/bin/run.sh
user={{ user_app }}
environment=JAVA_HOME="{{ jdk_dir }}",PATH="{{ jdk_dir }}/bin:%(ENV_PATH)s"
numprocs=1
umask=022
autostart=true
autorestart=true
startsecs=15
stopwaitsecs=15
stopsignal=TERM
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=syslog
